# nfc-web-interface.service
[Unit]
Description=NFC Tag Management Web Interface
After=multi-user.target network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/nfc-display
Environment="PATH=/home/pi/nfc-display/venv/bin"
ExecStart=/home/pi/nfc-display/venv/bin/python /home/pi/nfc-display/web_interface.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
